# netlogin.timer
[Unit]
Description=Run netlogin hourly

[Timer]
OnBootSec=2min
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
